# packer/scripts/maps/minecraft-map-refresh@.service
[Unit]
Description=Save world and rebuild uNmINeD map for %i
After=network.target minecraft@%i.service

[Service]
Type=oneshot
User=minecraft
EnvironmentFile=-/etc/minecraft.env
ExecStartPre=/usr/local/bin/mcrcon -H 127.0.0.1 -P ${RCON_PORT} -p ${RCON_PASSWORD} "say Saving world before map rebuild..."
ExecStartPre=/usr/local/bin/mcrcon -H 127.0.0.1 -P ${RCON_PORT} -p ${RCON_PASSWORD} save-all
ExecStart=/usr/local/bin/rebuild-map.sh /srv/minecraft-server/%i/world
