# packer/scripts/maps/minecraft-map-backup@.service
[Unit]
Description=Backup static maps for Minecraft world %i
After=network-online.target

[Service]
Type=oneshot
User=minecraft
EnvironmentFile=-/etc/minecraft.env
ExecStart=/usr/local/bin/backup-maps.sh %i
