# packer/scripts/minecraft/worlds/minecraft-world-prune.service
[Unit]
Description=Prune old Minecraft world backups
After=network-online.target

[Service]
Type=oneshot
User=minecraft
EnvironmentFile=-/etc/minecraft.env
ExecStart=/usr/local/bin/prune-worlds.sh
