# packer/scripts/minecraft/worlds/minecraft-world-prune.timer
[Unit]
Description=Weekly prune of world backups

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
