# packer/scripts/minecraft/worlds/minecraft-world-backup.timer
[Unit]
Description=Weekly full backup of all Minecraft worlds

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
Unit=minecraft-world-backup.service

[Install]
WantedBy=timers.target
