# packer/scripts/health/minecraft-health.service
[Unit]
Description=Minecraft stack health check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=minecraft
EnvironmentFile=-/etc/minecraft.env
ExecStart=/usr/local/bin/mc-healthcheck.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
