# packer/scripts/minecraft/maps/minecraft-map-build-daily@.timer
[Unit]
Description=Daily map build for %i (fallback when server not running)

[Timer]
OnCalendar=daily
Persistent=true
Unit=minecraft-map-build@%i.service

[Install]
WantedBy=timers.target
