# packer/scripts/minecraft/maps/minecraft-map-rebuild@.timer
[Unit]
Description=Rebuild Minecraft map for %i daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
