# packer/scripts/minecraft/maps/minecraft-map-refresh@.service
[Unit]
Description=Save world and rebuild uNmINeD map for %i
After=network.target minecraft@%i.service

[Service]
Type=oneshot
User=minecraft
EnvironmentFile=-/etc/minecraft.env
ExecStartPre=/usr/local/bin/minecraftctl rcon send "say Saving world before map rebuild..."
ExecStartPre=/usr/local/bin/minecraftctl rcon send "save-all"
ExecStart=/usr/local/bin/rebuild-map.sh /srv/minecraft-server/%i
