# packer/scripts/minecraft/maps/minecraft-map-update-daily@.timer
[Unit]
Description=Daily map update for %i (fallback when server not running)

[Timer]
OnCalendar=daily
Persistent=true
Unit=minecraft-map-update@%i.service

[Install]
WantedBy=timers.target
